# shariah-trader.service
[Unit]
Description=Shariah Algo Trader Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/shariah-algo-trader
EnvironmentFile=/home/ubuntu/shariah-algo-trader/.env
ExecStart=/home/ubuntu/shariah-algo-trader/.venv/bin/shariah-trader
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
